# example/hyprland-session.target
[Unit]
Description=Hyprland user session
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
