# systemd/hyprland-session.service.in
[Unit]
Description=Hyprland - Tiling compositor with the looks
Documentation=man:Hyprland(1)

BindsTo=hyprland-session.target
Before=hyprland-session.target

Wants=xdg-desktop-autostart.target

Wants=hyprland-session-pre.target
After=hyprland-session-pre.target

[Service]
Type=notify
ExecStart=@BINDIR@/Hyprland
ExecStop=@BINDIR@/hyprctl dispatch exit
ExecStopPost=systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY HYPRLAND_INSTANCE_SIGNATURE XDG_CURRENT_DESKTOP
OnFailure=systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY HYPRLAND_INSTANCE_SIGNATURE XDG_CURRENT_DESKTOP
OnSuccess=systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY HYPRLAND_INSTANCE_SIGNATURE XDG_CURRENT_DESKTOP
Restart=on-failure
Slice=session.slice
