# systemd/hyprland-session.service.in
[Unit]
Description=Hyprland - Tiling compositor with the looks
Documentation=man:Hyprland(1)
BindsTo=graphical-session.target
Before=graphical-session.target
Wants=xdg-desktop-autostart.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target

[Service]
Type=notify
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/Hyprland
ExecStop=@CMAKE_INSTALL_PREFIX@/bin/hyprctl dispatch exit
ExecStopPost=@CMAKE_INSTALL_PREFIX@/bin/systemctl --user unset-environment DISPLAY WAYLAND_DISPLAY HYPRLAND_INSTANCE_SIGNATURE XDG_CURRENT_DESKTOP
Restart=on-failure
Slice=session.slice
